# external/enarx-keep-nil@.service
[Unit]
Description=Enarx Keep-Loader (nil)
Documentation=https://github.com/enarx/

[Service]
Type=simple
Environment="ENARX_BACKEND=nil"
ExecStart=/home/mikecamel/programming/enarx-project/enarx-keepldr/target/debug/enarx-keepldr exec --sock %I /home/mikecamel/programming/enarx-project/enarx-wasmldr/target/x86_64-unknown-linux-musl/release/enarx-wasmldr
#/home/mike/programming/enarx-keepldr/test

UMask=0066

StandardOutput=journal+console
StandardError=journal+console

