# external/enarx-keep-sgx@.service
[Unit]
Description=Enarx Keep-Loader (sgx)
Documentation=https://github.com/enarx/

[Service]
Type=simple
Environment="ENARX_BACKEND=sgx"
ExecStart=/home/mike/programming/enarx-keepldr/target/debug/enarx-keepldr exec --sock I% ./test

UMask=0066

StandardOutput=journal+console
StandardError=journal+console

