# external/enarx-keep@.service
[Unit]
Description=Enarx Keep-Loader
Documentation=https://github.com/enarx/

[Service]
Type=simple
Environment="SCRIPT_ARGS=%I"
ExecStart=/home/mike/programming/enarx-keepldr/target/debug/enarx-keepldr $SCRIPT_ARGS

UMask=0066

StandardOutput=file:/home/mike/enarx-keepldr-%i.log
StandardError=file:/home/mike/enarx-keepldr-%i.err

