# external/enarx-keep-kvm@.service
[Unit]
Description=Enarx Keep-Loader (kvm)
Documentation=https://github.com/enarx/

[Service]
Type=simple
Environment="ENARX_BACKEND=kvm"
ExecStart=/home/mike/programming/enarx-project/enarx-keepldr/target/debug/enarx-keepldr exec --sock I% /home/mike/programming/enarx-project/enarx-wasmldr/target/x86_64-unknown-linux-musl/debug/enarx-wasmldr
#ExecStart=/home/mike/programming/enarx-keepldr/target/debug/enarx-keepldr exec --sock I% /home/mike/programming/enarx-keepldr/test

UMask=0066

StandardOutput=journal+console
StandardError=journal+console

